# deploy/jobs/crontab
0 6 * * * ./backup_pg.sh
0 5 * * * ./cleanup.sh
